# deploy/host/crictl-prune.timer
# Fires crictl-prune.service weekly. Sunday 04:00 KST — after the
# backend's 03:00 daily cron (cleanup + scan), before any meaningful
# user activity. Persistent=true so a missed window (host off) runs on
# next boot instead of waiting a full week.
[Unit]
Description=Weekly containerd image prune
Documentation=https://github.com/sungwookoo/swkoo-kr

[Timer]
OnCalendar=Sun *-*-* 04:00:00 Asia/Seoul
Persistent=true
Unit=crictl-prune.service

[Install]
WantedBy=timers.target
